# localai-git/localai.service
[Unit]
Description=LocalAI server

# "%C" Cache directory root: /var/cache (system) or "$XDG_CACHE_HOME" (user)
# "%E" Configuration directory root: /etc/ (system) or "$XDG_CONFIG_HOME" (user)
# "%S" State directory root: /var/lib (system) or $XDG_STATE_HOME (user)
# "%T" Directory for temporary files: /tmp or the path "$TMPDIR", "$TEMP" or "$TMP"
# "%N" Full unit name	Same as "%n", but with the type suffix removed

[Service]
User=localai
Group=localai
Type=exec
WorkingDirectory=%S/%N
Restart=on-failure
# PrivateTmp=yes

# default environment and local env
EnvironmentFile=%E/%N/%N.conf
EnvironmentFile=-%S/%N/.env

# start server
ExecStart=/usr/bin/localai run \
  --backend-assets-path="%S/%N" \
  --config-path="%S/%N" \
  --localai-config-dir="%S/%N/config" \
  --models-path="%S/%N/models" \
  --audio-path="%T/%N/audio" \
  --image-path="%T/%N/images" \
  --upload-path="%T/%N/upload" \
  --external-grpc-backends=bark:%S/%N/backend-assets/python/bark/run.sh,coqui:%S/%N/backend-assets/python/coqui/run.sh,diffusers:%S/%N/backend-assets/python/diffusers/run.sh,parler-tts:%S/%N/backend-assets/python/parler-tts/run.sh,rerankers:%S/%N/backend-assets/python/rerankers/run.sh,sentencetransformers:%S/%N/backend-assets/python/sentencetransformers/run.sh,transformers:%S/%N/backend-assets/python/transformers/run.sh,transformers-musicgen:%S/%N/backend-assets/python/transformers-musicgen/run.sh

# create virtualenvs for python backends, recreate if localai is newer than venv
# disabled python backends: autogptq exllama exllama2 mamba openvoice petals vall-e-x vllm
ExecStartPost=bash -c 'sleep 3;\
for i in bark coqui diffusers parler-tts rerankers sentencetransformers transformers transformers-musicgen; do \
  bedir="%S/%N/backend-assets/python/$i"; \
  if test /usr/bin/localai -nt $bedir/venv; then rm -r $bedir/venv; fi; \
  if test ! -d $bedir/venv; then \
    echo "re/creating $bedir"; \
    %S/%N/backend-assets/python/$i/install.sh; \
  fi; \
done'

# make some time for virtualenvs to be installed
TimeoutStartSec=120

[Install]
WantedBy=default.target
