# scripts/picoclaw.service
[Unit]
Description=PicoClaw Launcher web console service
Documentation=https://github.com/sipeed/picoclaw
After=network.target

[Service]
Type=simple
Environment=PICOCLAW_HOME=%h/.local/share/picoclaw
Environment=PICOCLAW_CONFIG=%h/.local/share/picoclaw/config.json
Environment=PICOCLAW_BINARY=/usr/bin/picoclaw
Environment=HOME=%h
EnvironmentFile=%h/.config/systemd/user/picoclaw.env
WorkingDirectory=%h
ExecStart=/usr/bin/picoclaw-launcher -no-browser
Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=picoclaw

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=

ProtectSystem=strict
TemporaryFileSystem=%h
BindPaths=%h/.local/share/picoclaw-service:%h/.local/share/picoclaw
BindPaths=%h/agent-shared:%h/agent-shared
ReadOnlyPaths=/etc/ssl /etc/ca-certificates /etc/resolv.conf /etc/hosts /etc/nsswitch.conf
PrivateTmp=yes
PrivateDevices=yes
PrivateMounts=yes

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
LockPersonality=yes

MemoryDenyWriteExecute=yes

SystemCallArchitectures=native

PrivateIPC=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
KeyringMode=private
UMask=0077

[Install]
WantedBy=default.target
